# wm/fluxbox/files/xdg/autostart/fluxbox.desktop
[Desktop Entry]
Exec=fluxbox 2> /var/log/fluxbox.log
